# www-servers/electroserver/files/5.2-electroserver.init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

opts="${opts} backup"

depend() {
	use dns logger net
}

start()	{
	ebegin "Starting ElectroServer 5.2"
	$ELECTROSERVER_HOME/ElectroServer_5_2_Service start
	eend $?
}

stop()	{
	ebegin "Stopping ElectroServer 5.2"
	$ELECTROSERVER_HOME/ElectroServer_5_2_Service stop
	eend $?
}

restart() {
	ebegin "Restarting ElectroServer 5.2"
  $ELECTROSERVER_HOME/ElectroServer_5_2_Service restart
	eend $?
}

status() {
	$ELECTROSERVER_HOME/ElectroServer_5_2_Service status
}

backup() {
  stop
  $ELECTROSERVER_HOME/DbUtility
  start
}
